# com.macfancontrol.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <!-- NOTE: Do not edit this file directly.
         Run install.sh to generate and install the plist with correct paths
         for your system. This file is a reference template only. 
         Manual load/unload: sudo launchctl load/unload /Library/LaunchAgents/com.macfancontrol.daemon.plist -->

    <key>Label</key>
    <string>com.macfancontrol.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/YOUR/INSTALL/PATH/venv/bin/python3</string>
        <string>/YOUR/INSTALL/PATH/daemon.py</string>
        <string>--config</string>
        <string>/YOUR/INSTALL/PATH/config.json</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/YOUR/INSTALL/PATH/logs/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/YOUR/INSTALL/PATH/logs/daemon.err</string>

    <key>WorkingDirectory</key>
    <string>/YOUR/INSTALL/PATH</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

</dict>
</plist>
